# etc/sv/ampr-ripd/run
#!/bin/sh

exec ampr-ripd -d -i tun44 -t 44 2>&1
