# etc/sv/socat-gopher/run
#!/bin/sh

exec socat tcp-listen:70,fork,reuseaddr tcp:swich:70
