# etc/sv/tunnel-blueiris/run
#!/bin/sh

exec autossh -M 0 -o ExitOnForwardFailure=yes -o TCPKeepAlive=no -o ServerAliveInterval=5 -o ServerAliveCountMax=6 trench -R 0.0.0.0:60000:win:8000 -N -i ~/.ssh/id_tunnel_blueiris
