# etc/sv/wpa_supplicant_atnt/run
#!/bin/sh

exec wpa_supplicant -Dwired -ieno1 -c/root/8021x/wpa_supplicant.conf 2>&1
